# config/nvim-lsp.service
[Unit]
Description=Neovim LSP Docker container
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=__SCRIPT_PATH__
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
